# deploy/social-economy/com.etzhayyim.kotoba.social-economy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the Mishmar social-capital economy tick (ADR-2606082100),
  to run on a Murakumo Mac mini fleet node (it can reach the geth-private tunnel
  + the KaizenObserver feed). One tick per StartInterval (default 1 day = the
  social-capital epoch); launchd owns the cadence, the binary is one-shot.

  Install (per README.md):
    cp com.etzhayyim.kotoba.social-economy.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.etzhayyim.kotoba.social-economy.plist
  Edit the run-tick.sh path below to where this deploy/ dir lives on the node.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.etzhayyim.kotoba.social-economy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/Shared/kotoba/deploy/social-economy/run-tick.sh</string>
  </array>

  <!-- one epoch = one day; launchd fires the one-shot tick on this cadence. -->
  <key>StartInterval</key>
  <integer>86400</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/Shared/kotoba/social-economy/tick.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/Shared/kotoba/social-economy/tick.err.log</string>

  <!-- read+verify only; no inbound listener. -->
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
