# crates/kotoba-kotodama/py/src/kotodama/projects/ameno/ameno-daemon.service
# systemd unit for ameno-daemon Path B (Python).
#
# Install (Linux):
#   sudo cp ameno-daemon.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ameno-daemon
#   journalctl -u ameno-daemon -f
#
# Edit User=, Environment paths, and ExecStart to match your host.
# Assumes Ollama is already running on localhost:11434.
#
# Authoritative ADR: 90-docs/adr/2605191257-ameno-daemon-path-b-kotodama-python.md

[Unit]
Description=Ameno headless daemon (Path B, kotodama Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/etzhayyim-root
Environment=AMENO_HOME=/home/YOUR_USERNAME/.ameno
Environment=AMENO_HOST=127.0.0.1
Environment=AMENO_PORT=12481
Environment=LOCAL_LLM_PROVIDER=ollama
Environment=LOCAL_LLM_ENDPOINT=http://127.0.0.1:11434/api/chat
Environment=LOCAL_LLM_MODEL=gemma3:4b
# uv-managed venv path; replace if you use poetry / pip / pipx.
ExecStart=/home/YOUR_USERNAME/etzhayyim-root/40-engine/kotoba/crates/kotoba-kotodama/py/.venv/bin/python -m kotodama.projects.ameno
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
